# build/supervisor/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php5-fpm -F
autostart=true
autorestart=true
redirect_stderr=true